# python-socks-proxy/debian/python-socks-proxy.service
#
# This file is managed by a debian package. Do not edit.
#

[Unit]
Description=Python Socks Proxy
After=syslog.target network.target

[Service]
User=root
Group=root
ExecStart=/usr/bin/python3 /usr/bin/server.py
Restart=always
RestartSec=10
KillSignal=SIGINT
StartLimitBurst=5
StartLimitInterval=0

[Install]
WantedBy=multi-user.target
